# screen_grabber/ronaldo.target
[Unit]
Description=Ronaldo backend services
Wants=xvfb.service vnc_server.service rtmp_server.service
